# docs/systemd/gws-daemon.service
[Unit]
Description=gws TUI daemon
After=network-online.target

[Service]
Type=simple
ExecStart=%h/go/bin/gws daemon start
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
